# backup.cron
3 * * * * /usr/local/bin/backup create --full
*/10 * * * * /usr/local/bin/backup create
0 * * * * /usr/local/bin/backup prune --retention "7:days,4:weeks,6:months,5:years"
